# kulgadd.service.in
[Unit]
Description=KU LGAD daemon
After=network.target

[Service]
ExecStart=@EXEC_START_PATH@
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
